# roles/bcm_prepare/files/fix-dual-nic-arp.service
[Unit]
Description=Fix ARP for dual NICs on same subnet
After=network-online.target cmfirstboot.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/fix-dual-nic-arp.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
